# gridpane-monitor.service
[Unit]
Description=GridPane Server Monitoring Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=gridpane-monitor
Group=gridpane-monitor
WorkingDirectory=/opt/gridpane-monitor
ExecStart=/usr/bin/python3 /opt/gridpane-monitor/gridpane_monitor.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gridpane-monitor
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=1024
LimitNPROC=64
MemoryMax=128M
CPUQuota=50%

# Network security
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=127.0.0.1

[Install]
WantedBy=multi-user.target
